# deploy-elk-stack-app/kanban-ui/supervisord-nginx-logstash.ini
[program:logstash]
command=/usr/share/logstash-7.9.0/bin/logstash -f /test.conf

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true
